# .github/deployment/production/supervisor-workers.conf
; ACME Corp Production Queue Workers Configuration
; Manages all Laravel queue workers for the CSR platform production environment

[group:acme-workers]
programs=acme-default,acme-payments,acme-notifications,acme-exports,acme-reports,acme-maintenance,acme-bulk,acme-cache-warming,acme-scout-indexer
priority=999

; Default Queue Worker (general tasks) - Production
[program:acme-default]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/acme-corp/current/artisan queue:work redis --sleep=3 --tries=3 --timeout=300 --memory=512
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/acme-corp/storage/logs/worker-default.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stopwaitsecs=3600
startsecs=0

; Payments Queue Worker (high priority, critical) - Production
[program:acme-payments]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/acme-corp/current/artisan queue:work redis --queue=payments --sleep=1 --tries=5 --timeout=180 --memory=256
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=3
redirect_stderr=true
stdout_logfile=/var/www/acme-corp/storage/logs/worker-payments.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stopwaitsecs=180
startsecs=0
priority=100

; Notifications Queue Worker (high priority) - Production
[program:acme-notifications]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/acme-corp/current/artisan queue:work redis --queue=notifications --sleep=2 --tries=3 --timeout=120 --memory=128
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/acme-corp/storage/logs/worker-notifications.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stopwaitsecs=120
startsecs=0
priority=200

; Exports Queue Worker (memory intensive) - Production
[program:acme-exports]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/acme-corp/current/artisan queue:work redis --queue=exports --sleep=10 --tries=2 --timeout=900 --memory=1024
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/acme-corp/storage/logs/worker-exports.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stopwaitsecs=900
startsecs=0
priority=500

; Reports Queue Worker (medium priority) - Production
[program:acme-reports]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/acme-corp/current/artisan queue:work redis --queue=reports --sleep=10 --tries=3 --timeout=600 --memory=256
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/acme-corp/storage/logs/worker-reports.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stopwaitsecs=600
startsecs=0
priority=400

; Maintenance Queue Worker (low priority) - Production
[program:acme-maintenance]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/acme-corp/current/artisan queue:work redis --queue=maintenance --sleep=15 --tries=1 --timeout=1800 --memory=256
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/acme-corp/storage/logs/worker-maintenance.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stopwaitsecs=1800
startsecs=0
priority=900

; Bulk Queue Worker (low priority, memory intensive) - Production
[program:acme-bulk]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/acme-corp/current/artisan queue:work redis --queue=bulk --sleep=20 --tries=1 --timeout=3600 --memory=1024
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=1
redirect_stderr=true
stdout_logfile=/var/www/acme-corp/storage/logs/worker-bulk.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stopwaitsecs=3600
startsecs=0
priority=999

; Cache Warming Queue Worker - Production
[program:acme-cache-warming]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/acme-corp/current/artisan queue:work redis --queue=cache-warming,cache-warming-orchestrator --sleep=5 --tries=3 --timeout=900 --memory=256
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=2
redirect_stderr=true
stdout_logfile=/var/www/acme-corp/storage/logs/worker-cache-warming.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stopwaitsecs=900
startsecs=0
priority=300

; Scout Indexer Queue Worker (Meilisearch search indexing) - Production
[program:acme-scout-indexer]
process_name=%(program_name)s_%(process_num)02d
command=/usr/bin/php /var/www/acme-corp/current/artisan queue:work redis --queue=scout --sleep=1 --tries=3 --timeout=300 --memory=512 --max-jobs=2000
autostart=true
autorestart=true
stopasgroup=true
killasgroup=true
user=www-data
numprocs=4
redirect_stderr=true
stdout_logfile=/var/www/acme-corp/storage/logs/worker-scout-indexer.log
stdout_logfile_maxbytes=10MB
stdout_logfile_backups=10
stopwaitsecs=300
startsecs=0
priority=150